# etc/stub-runner.jar
#!/bin/bash
#
#    .   ____          _            __ _ _
#   /\\ / ___'_ __ _ _(_)_ __  __ _ \ \ \ \
#  ( ( )\___ | '_ | '_| | '_ \/ _` | \ \ \ \
#   \\/  ___)| |_)| | | | | || (_| |  ) ) ) )
#    '  |____| .__|_| |_|_| |_\__, | / / / /
#   =========|_|==============|___/=/_/_/_/
#   :: Spring Boot Startup Script ::
#

### BEGIN INIT INFO
# Provides:          spring-cloud-contract-stub-runner-boot
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Spring Cloud Contract Stub Runner Boot
# Description:       Spring Cloud Contract Stub Runner Boot
# chkconfig:         2345 99 01
### END INIT INFO

[[ -n "$DEBUG" ]] && set -x

# Initialize variables that cannot be provided by a .conf file
WORKING_DIR="$(pwd)"
# shellcheck disable=SC2153
[[ -n "$JARFILE" ]] && jarfile="$JARFILE"
[[ -n "$APP_NAME" ]] && identity="$APP_NAME"

# Follow symlinks to find the real jar and detect init.d script
cd "$(dirname "$0")" || exit 1
[[ -z "$jarfile" ]] && jarfile=$(pwd)/$(basename "$0")
while [[ -L "$jarfile" ]]; do
  if [[ "$jarfile" =~ init\.d ]]; then
    init_script=$(basename "$jarfile")
  else
    configfile="${jarfile%.*}.conf"
    # shellcheck source=/dev/null
    [[ -r ${configfile} ]] && source "${configfile}"
  fi
  jarfile=$(readlink "$jarfile")
  cd "$(dirname "$jarfile")" || exit 1
  jarfile=$(pwd)/$(basename "$jarfile")
done
jarfolder="$( (cd "$(dirname "$jarfile")" && pwd -P) )"
cd "$WORKING_DIR" || exit 1

# Inline script specified in build properties


# Source any config file
configfile="$(basename "${jarfile%.*}.conf")"

# Initialize CONF_FOLDER location defaulting to jarfolder
[[ -z "$CONF_FOLDER" ]] && CONF_FOLDER="${jarfolder}"
# shellcheck source=/dev/null
[[ -r "${CONF_FOLDER}/${configfile}" ]] && source "${CONF_FOLDER}/${configfile}"

# ANSI Colors
echoRed() { echo $'\e[0;31m'"$1"$'\e[0m'; }
echoGreen() { echo $'\e[0;32m'"$1"$'\e[0m'; }
echoYellow() { echo $'\e[0;33m'"$1"$'\e[0m'; }

# Initialize PID/LOG locations if they weren't provided by the config file
[[ -z "$PID_FOLDER" ]] && PID_FOLDER="/var/run"
[[ -z "$LOG_FOLDER" ]] && LOG_FOLDER="/var/log"
! [[ "$PID_FOLDER" == /* ]] && PID_FOLDER="$(dirname "$jarfile")"/"$PID_FOLDER"
! [[ "$LOG_FOLDER" == /* ]] && LOG_FOLDER="$(dirname "$jarfile")"/"$LOG_FOLDER"
! [[ -x "$PID_FOLDER" ]] && echoYellow "PID_FOLDER $PID_FOLDER does not exist. Falling back to /tmp" && PID_FOLDER="/tmp"
! [[ -x "$LOG_FOLDER" ]] && echoYellow "LOG_FOLDER $LOG_FOLDER does not exist. Falling back to /tmp" && LOG_FOLDER="/tmp"

# Set up defaults
[[ -z "$MODE" ]] && MODE="auto" # modes are "auto", "service" or "run"
[[ -z "$USE_START_STOP_DAEMON" ]] && USE_START_STOP_DAEMON="true"

# Create an identity for log/pid files
if [[ -z "$identity" ]]; then
  if [[ -n "$init_script" ]]; then
    identity="${init_script}"
  else
    identity=$(basename "${jarfile%.*}")_${jarfolder//\//}
  fi
fi

# Initialize log file name if not provided by the config file
[[ -z "$LOG_FILENAME" ]] && LOG_FILENAME="${identity}.log"

# Initialize stop wait time if not provided by the config file
[[ -z "$STOP_WAIT_TIME" ]] && STOP_WAIT_TIME="60"

# Utility functions
checkPermissions() {
  touch "$pid_file" &> /dev/null || { echoRed "Operation not permitted (cannot access pid file)"; return 4; }
  touch "$log_file" &> /dev/null || { echoRed "Operation not permitted (cannot access log file)"; return 4; }
}

isRunning() {
  ps -p "$1" &> /dev/null
}

await_file() {
  end=$(date +%s)
  let "end+=10"
  while [[ ! -s "$1" ]]
  do
    now=$(date +%s)
    if [[ $now -ge $end ]]; then
      break
    fi
    sleep 1
  done
}

# Determine the script mode
action="run"
if [[ "$MODE" == "auto" && -n "$init_script" ]] || [[ "$MODE" == "service" ]]; then
  action="$1"
  shift
fi

# Build the pid and log filenames
PID_FOLDER="$PID_FOLDER/${identity}"
pid_file="$PID_FOLDER/${identity}.pid"
log_file="$LOG_FOLDER/$LOG_FILENAME"

# Determine the user to run as if we are root
# shellcheck disable=SC2012
[[ $(id -u) == "0" ]] && run_user=$(ls -ld "$jarfile" | awk '{print $3}')

# Ensure the user actually exists
id -u "$run_user" &> /dev/null || unset run_user

# Run as user specified in RUN_AS_USER
if [[ -n "$RUN_AS_USER" ]]; then
    if ! [[ "$action" =~ ^(status|run)$ ]]; then
        id -u "$RUN_AS_USER" || {
            echoRed "Cannot run as '$RUN_AS_USER': no such user"
            exit 2
        }
        [[ $(id -u) == 0 ]] || {
            echoRed "Cannot run as '$RUN_AS_USER': current user is not root"
            exit 4
        }
    fi
    run_user="$RUN_AS_USER"
fi

# Issue a warning if the application will run as root
[[ $(id -u ${run_user}) == "0" ]] && { echoYellow "Application is running as root (UID 0). This is considered insecure."; }

# Find Java
if [[ -n "$JAVA_HOME" ]] && [[ -x "$JAVA_HOME/bin/java" ]]; then
    javaexe="$JAVA_HOME/bin/java"
elif type -p java > /dev/null 2>&1; then
    javaexe=$(type -p java)
elif [[ -x "/usr/bin/java" ]];  then
    javaexe="/usr/bin/java"
else
    echo "Unable to find Java"
    exit 1
fi

arguments=(-Dsun.misc.URLClassPath.disableJarChecking=true $JAVA_OPTS -jar "$jarfile" $RUN_ARGS "$@")

# Action functions
start() {
  if [[ -f "$pid_file" ]]; then
    pid=$(cat "$pid_file")
    isRunning "$pid" && { echoYellow "Already running [$pid]"; return 0; }
  fi
  do_start "$@"
}

do_start() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  if [[ ! -e "$PID_FOLDER" ]]; then
    mkdir -p "$PID_FOLDER" &> /dev/null
    if [[ -n "$run_user" ]]; then
      chown "$run_user" "$PID_FOLDER"
    fi
  fi
  if [[ ! -e "$log_file" ]]; then
    touch "$log_file" &> /dev/null
    if [[ -n "$run_user" ]]; then
      chown "$run_user" "$log_file"
    fi
  fi
  if [[ -n "$run_user" ]]; then
    checkPermissions || return $?
    if [ $USE_START_STOP_DAEMON = true ] && type start-stop-daemon > /dev/null 2>&1; then
      start-stop-daemon --start --quiet \
        --chuid "$run_user" \
        --name "$identity" \
        --make-pidfile --pidfile "$pid_file" \
        --background --no-close \
        --startas "$javaexe" \
        --chdir "$working_dir" \
        -- "${arguments[@]}" \
        >> "$log_file" 2>&1
      await_file "$pid_file"
    else
      su -s /bin/sh -c "$javaexe $(printf "\"%s\" " "${arguments[@]}") >> \"$log_file\" 2>&1 & echo \$!" "$run_user" > "$pid_file"
    fi
    pid=$(cat "$pid_file")
  else
    checkPermissions || return $?
    "$javaexe" "${arguments[@]}" >> "$log_file" 2>&1 &
    pid=$!
    disown $pid
    echo "$pid" > "$pid_file"
  fi
  [[ -z $pid ]] && { echoRed "Failed to start"; return 1; }
  echoGreen "Started [$pid]"
}

stop() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  [[ -f $pid_file ]] || { echoYellow "Not running (pidfile not found)"; return 0; }
  pid=$(cat "$pid_file")
  isRunning "$pid" || { echoYellow "Not running (process ${pid}). Removing stale pid file."; rm -f "$pid_file"; return 0; }
  do_stop "$pid" "$pid_file"
}

do_stop() {
  kill "$1" &> /dev/null || { echoRed "Unable to kill process $1"; return 1; }
  for ((i = 1; i <= STOP_WAIT_TIME; i++)); do
    isRunning "$1" || { echoGreen "Stopped [$1]"; rm -f "$2"; return 0; }
    [[ $i -eq STOP_WAIT_TIME/2 ]] && kill "$1" &> /dev/null
    sleep 1
  done
  echoRed "Unable to kill process $1";
  return 1;
}

force_stop() {
  [[ -f $pid_file ]] || { echoYellow "Not running (pidfile not found)"; return 0; }
  pid=$(cat "$pid_file")
  isRunning "$pid" || { echoYellow "Not running (process ${pid}). Removing stale pid file."; rm -f "$pid_file"; return 0; }
  do_force_stop "$pid" "$pid_file"
}

do_force_stop() {
  kill -9 "$1" &> /dev/null || { echoRed "Unable to kill process $1"; return 1; }
  for ((i = 1; i <= STOP_WAIT_TIME; i++)); do
    isRunning "$1" || { echoGreen "Stopped [$1]"; rm -f "$2"; return 0; }
    [[ $i -eq STOP_WAIT_TIME/2 ]] && kill -9 "$1" &> /dev/null
    sleep 1
  done
  echoRed "Unable to kill process $1";
  return 1;
}

restart() {
  stop && start
}

force_reload() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  [[ -f $pid_file ]] || { echoRed "Not running (pidfile not found)"; return 7; }
  pid=$(cat "$pid_file")
  rm -f "$pid_file"
  isRunning "$pid" || { echoRed "Not running (process ${pid} not found)"; return 7; }
  do_stop "$pid" "$pid_file"
  do_start
}

status() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  [[ -f "$pid_file" ]] || { echoRed "Not running"; return 3; }
  pid=$(cat "$pid_file")
  isRunning "$pid" || { echoRed "Not running (process ${pid} not found)"; return 1; }
  echoGreen "Running [$pid]"
  return 0
}

run() {
  pushd "$(dirname "$jarfile")" > /dev/null
  "$javaexe" "${arguments[@]}"
  result=$?
  popd > /dev/null
  return "$result"
}

# Call the appropriate action function
case "$action" in
start)
  start "$@"; exit $?;;
stop)
  stop "$@"; exit $?;;
force-stop)
  force_stop "$@"; exit $?;;
restart)
  restart "$@"; exit $?;;
force-reload)
  force_reload "$@"; exit $?;;
status)
  status "$@"; exit $?;;
run)
  run "$@"; exit $?;;
*)
  echo "Usage: $0 {start|stop|force-stop|restart|force-reload|status|run}"; exit 1;
esac

exit 0
PK  &U            	  META-INF/   PK           PK  &U               META-INF/MANIFEST.MF}j0_f9&c)Yԑ4ml$dӯ1&E.P
$0z7.y,U;^>LVՁCTUB~6yMIb,Uw09ttWYȧ\fW\vs1å&^ۣnGtoV5#$.U'$V~?γEoGelrf\7GYRHPKh     PK
    SR               org/PK
    SR               org/springframework/PK
    SR               org/springframework/boot/PK
    SR                org/springframework/boot/loader/PK
    SR            (   org/springframework/boot/loader/wrapper/PK
    XTR               META-INF/maven/PK
    XTR            5   META-INF/maven/org.springframework.boot.experimental/PK
    XTR            N   META-INF/maven/org.springframework.boot.experimental/spring-boot-thin-wrapper/PK  SR            <   org/springframework/boot/loader/wrapper/ThinJarWrapper.classZ	xT>g{\	A*dCP""$7d&ZjvBjEkebՊպԮZVkV%G{޿(iN4]"]`#V&j`Ṯ'YZC\Lc5`fuV#TnE}ӖWc6dכa>Ө%H<$ISItu<u~ڵm:QP{,ەPvɊu:Io:zB&Sмe{1dXfN%f_Yk[[ө3T,[ܲ\IN=u4֮hoiXAp#fnYmKú-Kju:أz=Ֆ=GdHUwSIi;aPIkG,T	}[]=zykjpX!3XvN&B5ym+70zmXK2+$iL&?gA(J,Cfx9,S$PfY(.6{I8>	->	-ZVfSX(ft>ji4h1-HlCxU-c:d{ZA HKcD#femFjj0m"ڲ~ızCq2+V^VBmGɄ&a_2t!mШmHiOmMH<sCA@0tA[,ٔm7caS8`,$c{uIƏqðe͎H/-JnN;5ԶBrRX,31E:*zl~/;K'swW0~کюTtEC$v2J)7KAeP4DFn((ӕL*IO! ֐@$_X ۯ5閭PH,|Π2:!d I/"}D}٠"A%3ͣ`9*TN{H.wFq%Ul|A	3#"]$*R|Ahmg(+EtZD}8lvYDpb:e8PۈZKH4Q*2afRQi!׽Gw~Cx|!qРw"j-ҝt{"a#C43J`y,qn|}?2Aѩ&GD?4a:4Z»xfqOf~dУFQ
Fy$=44әbA?APLDڟK!fvk"L<EZg9/+ 3Ht<fc{<#OGN9t@MN'F@jjhПEM'kF4#j+blk$֢Č. T#A]$nBJf0MzsRV%{NAoӿ1cY,CY;. NA@Fa)gu
`z?RuDPa .N]2`GBܒ=k[Қ`u)OBTul8+h++`BA ð~Reb1FBJ2]bC:8'0qjjī`_'y"Oʺ>o.) Bg&ZVHg4V>7XKnPaČ'N+7hs{*b]H|&1^UUn'W\q%R8FNi<IGI):c7<:#fFr*fv`,]=P+X^dp-/Qu-|" 0vReZ3-.i5^j#\jQxR"/oFB]fvX[ɹBn5T`(R8(9GPvpC_(G,B[ӡ9b`M@.MR;ѸE9Y~uNCL*qɔ$(ݖWcsH﹮'!%CʛhsJ.9419"JPryu?DCqО4kи?]-_N2vnaxofDvd"r!'RԖ|SyڃI(O_1{U`W'Ņ1oM"BS_OGΡuHpI2fPCIt"%+Q2*I&Ȋh[܍=Q_1ZnTkPh2nƪ{Qt\-=lg|&HVMpP,=i9q2#o[/IpsI0B>bf rn3q*ܪ ay:'{`$	!)Zm; :%P<(h:N(Ö!T?ϖ91!vu[v֫Y7 GĪ:9%G}/4a:9uu2Y.\gL'xcQ0K	ѰS?I782wX(Eh8A(G7+]D4]I	!˹' U2T*ql\,1vs}~P:lU@@~g*/??3(!T/	31O[ɑ
*>LXyfvA:J5~=ߐf@&Ӵwm;)_v3SAl.)Ku~ime|H(y4/]&;b1 ӌA[
u͛͘Z)HuwÉ[Q%@10nڤ;Ps~24ޛo=:io,IYI,3b!+>>ZBkNG$"|&FF1.a<*g<&؜8c<!g|
ƾD79g<rƥOOϙtzAex*'?vr^n8BJzJМM-MAZCZU>Ӡ"Ժn*#.=4y vZws Q9_4Zo[1SDh?r}($gEvw"<Nc%G.Aؽ=f$
QQCMAqEMD2pPCUNеRitVO>!SMBj	瀥ت]6HhG縜Q=epEMs!54ҝb})yRtsԵF<.#wwzn;1X1e$~];p4vS+>E]"m!6(,h!ݢ.>r <겗x!:AHTqB'#WXmp+pу}Q[9HH7Ф#tg:{&RA<.}l7Cjjq&HUs豀sO=r9>W 9Ejͺ&&0uPr%$5k4m@՟U^hAan y&-t=9(LT.85RT8h8h-"FCq"0n6M\5TxFkj@a7AiT{SQB5KM2#˕lJk\wRnrW xb*C:PCOSdz-^`MI6sx|*eC(yK{(O;NOGb>EW^iP3žhZ
'(,Y] ##}1:e{_:C+kM7`t+*&gh?Nqt7O{x>^F.oAg~t0GN:܇|pPOU9i%])8(K<y>p`C9@YSk|!!UdYgMSpw) A?tQFsޣ^o)$ƹv
<'W<B!_̀~XxB?#H]~E'˓oUyN`e>mg4Us<$`Jܳ\@:~^w/ 5)Mk*e˙TsЏќs<BT/
}b>GD݇ypDQC<0/`}t ~	I)	s?5Ī_!&=\kC|P}nVԷ?*z(խ<@[T[CuTm_NOU$$2麜ܹK*>s rMSJ;*>n:yESU-4Jj,1|;U9R;2s}Mx2C[r ʝOG6=H:{1Szn08g 5AexpOa+.
ͦWApzޢ/Ap{F`܁q:p{BNN2+m NF6Up 1%dYkL&1UfeenR;ڐUΓX_;bC*ȴ;(歨9 ۲B9d\F9ukFbpDI`͙Mq_sNUz9;R| _]]50-x<;h;iC)gf$ijMO#HQ@/gUJV X% Z@T?GhdB:B&ꄖrra1hs,{Ŏ8]1vv69F)<]x<}v19YK//Oߠa&aeT
 ߘh=OAmX|fLlU蔯6ËlՍZ|bKN*,=ЋgoeaHAΔ,~NGEoAdKT⯣Ԓ b5@a~ĕ4 ݚf?M͓pi-f"^s3iC/opvNmۛ=.Q'g
>Z<&)Plr
U:LĶ|^*9[+
,.E}@;qH:\Rަ[2;d׮dDMpA3M1Ж׸pH4)|ng'O'K<7xNp?(IPC]ȆI!8}}\O;ZzWh/s:܋J  oV.Gy-χQ<~OL7&^S\Bg#r+h
^Ixpm{;C^O[v&r]Ϋ
^C^ۏ}wc}7}Ǹ~(֞ϰ֞GJ}
*kǁ2PyɆ=RJ3,2T	 oyx-\/ډZN=dR)VkTKmQi+{iZYR$2,0~_:ĿoDiCxa{[6_v>V
-JóXV5VzO&4_4  ~Ӈ;(=)J	nM \@lB7]Hm=pCy8LKm(ur	<,e4&,{p8UYC }QuG-]TYl:ÌհuRC>
=vHIj"njH0PK  .  PK  &U               lib/ PK           PK  &U            
   lib/.empty PK           PK
    &U               org/springframework/cloud/PK
    &U            #   org/springframework/cloud/contract/PK
    &U            .   org/springframework/cloud/contract/stubrunner/PK
    &U            5   org/springframework/cloud/contract/stubrunner/server/PK
    &U            )   META-INF/maven/org.springframework.cloud/PK
    &U            P   META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/PK  &U            I   org/springframework/cloud/contract/stubrunner/server/StubRunnerBoot.classRn1}N$ZZ("n!|A,Ԯ)pݸ[y	ā࣪H	R |ьg{۟ A:vx1xxrGk#zbƒak|WMRi?Lg`2QU>'Rcs^N/<+L5ΊU^RCȜbv{WZB蜪<l6qg]{pRa44[w9`أS9R"JZ'2ڏVJbvJ>%x:-T hT{uC:nc?h~OT[(y^1zs#Ð^$)DYƫpq</Hf1^Fןl"N_[d!Y"lߠU :}ZGz}lx_ w;PKyZa  v  PK  &U               META-INF/thin-rabbit.properties  L0 (A
)_5y]wNMkA!E;•3v+
ܪdsUȞ$;VPK_X      PK  &U               META-INF/thin-kafka.properties	  Nt %XDEnσ;ONM{A!E;•a*Ѣ@J6˷
f1jPK]X      PK  &U               application.yml- ;OpgѨpt'	cd_`;69qďU .-z' B2sATܶX,<RW:/łνkFe1PnLff"ѕ-!{=ucs0<W?ɻPKj6      PK  f&U            W   META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.xmlUQo0~"JPHULDt^Mr`GT.vQ-|Ηv5(ͥ|Loϟm܋
%_ 1>n5c!͆lTsz菇4Y\hDAqGZ5
HK}@oH?oM2	3VY_!1V.Z sh%Sȿc*{1L0Es%>1%хe`#Ւ$,ӈ6
ΔK:0(.uD[W_zŋ؍5<2	^ ҝצT8PV,SUOmZJi)lZMkHJA'o9&%@m?̂ՑI8=yOv G$ 7IwxA1:C@U7QKir.7YgT]60;uM<w"/qH;KU|YЌKꑵABR6!6F:u{lq
<]:9{egEYZs<d `juڹz_[Tqږ1g_3ht]NlGֳjPO{PKIZjh    PK  &U            ^   META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.properties%K E9{))p/Dmͣ%:CkAڅʅ!=RVNMF	)~]PO2n?kVl'PK[ST   b   PK  &U           	         A9$  META-INF/  PK  &Uh                v$  META-INF/MANIFEST.MFPK
    SR                      A%  org/PK
    SR                      A%  org/springframework/PK
    SR                      A%  org/springframework/boot/PK
    SR                       A5&  org/springframework/boot/loader/PK
    SR            (          As&  org/springframework/boot/loader/wrapper/PK
    XTR                      A&  META-INF/maven/PK
    XTR            5          A&  META-INF/maven/org.springframework.boot.experimental/PK
    XTR            N          A9'  META-INF/maven/org.springframework.boot.experimental/spring-boot-thin-wrapper/PK  SR  .  <           '  org/springframework/boot/loader/wrapper/ThinJarWrapper.classPK  &U                     A?  lib/PK  &U           
           /@  lib/.emptyPK
    &U                      Ai@  org/springframework/cloud/PK
    &U            #          A@  org/springframework/cloud/contract/PK
    &U            .          A@  org/springframework/cloud/contract/stubrunner/PK
    &U            5          A.A  org/springframework/cloud/contract/stubrunner/server/PK
    &U            )          AA  META-INF/maven/org.springframework.cloud/PK
    &U            P          AA  META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/PK  &UyZa  v  I           6B  org/springframework/cloud/contract/stubrunner/server/StubRunnerBoot.classPK  &U_X                 oD  META-INF/thin-rabbit.propertiesPK  &U]X                 E  META-INF/thin-kafka.propertiesPK  &Uj6                 E  application.ymlPK  f&UIZjh    W           F  META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.xmlPK  &U[ST   b   ^           nI  META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.propertiesPK      W  NJ    